# crontab
*  *  *  *  * /bin/bash -l /opt/ArDNSPod-IPv6/ddnspod-ipv6.sh >> /tmp/ddns.log
